# templates/default/snowflake.erb
#!/bin/sh
#
# snowflake    Init script for snowflake
#
# chkconfig: 345 99 01
# processname: snowflake
# pidfile: var/run/snowflake.pid
# description: Starts and Stops and restarts snowflake process.
#

java -server -XX:+UseConcMarkSweepGC -verbosegc -XX:+PrintGCDetails -XX:+PrintGCTimeStamps -XX:+PrintGCDateStamps -XX:+UseParNewGC -Xloggc:/var/log/snowflake/gc.log -Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=9999 -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -Xmx700m -Xms700m -Xmn500m -XX:ErrorFile=/var/log/snowflake/java_error%p.log -cp /usr/local/snowflake/target/lib/*:/usr/local/snowflake/target/snowflake-1.0.1-SNAPSHOT.jar com.twitter.service.snowflake.SnowflakeServer -f /usr/local/snowflake/config/config.scala

APP_NAME="<%= node['snowflake-nativex']['app']['nativex_snowflake_project_name'] %>"
APP_HOME="<%= node['snowflake-nativex']['app']['snowflake_home'] %>"

DEBUG_OPTS="<%= node['snowflake-nativex']['java']['debug_opts'] %>"
GC_OPTS="<%= node['snowflake-nativex']['java']['gc_opts'] %>"
HEAP_OPTS="<%= node['snowflake-nativex']['java']['heap_opts'] %>"
JAVA_OPTS="<%= node['snowflake-nativex']['java']['java_opts'] %>"
JMX_OPTS="<%= node['snowflake-nativex']['java']['jmx_opts'] %>"
MAIN_CLASS="<%= node['snowflake-nativex']['java']['main-class'] %>"
MAIN_JAR="<%= node['snowflake-nativex']['java']['main_jar'] %>"

CLASSPATH_OPT="$APP_HOME/target/lib/*:$APP_HOME/target/$MAIN_JAR"

PID_PATH_NAME="/var/run/$APP_NAME.pid"

case $1 in
    start)
        echo "Starting $APP_NAME ..."
        if [ ! -f $PID_PATH_NAME ]; then
            nohup java ${JAVA_OPTS} -cp ${CLASSPATH_OPT} ${MAIN_CLASS} -f ${APP_HOME}/config/config.scala /tmp 2>> /dev/null >> /dev/null &
                        echo $! > $PID_PATH_NAME
            echo "$APP_NAME started ..."
        else
            echo "$APP_NAME is already running ..."
        fi
    ;;
    stop)
        if [ -f $PID_PATH_NAME ]; then
            PID=$(cat $PID_PATH_NAME);
            echo "$APP_NAME stoping ..."
            kill $PID;
            echo "$APP_NAME stopped ..."
            rm $PID_PATH_NAME
        else
            echo "$APP_NAME is not running ..."
        fi
    ;;
    restart)
        if [ -f $PID_PATH_NAME ]; then
            PID=$(cat $PID_PATH_NAME);
            echo "$APP_NAME stopping ...";
            kill $PID;
            echo "$APP_NAME stopped ...";
            rm $PID_PATH_NAME
            echo "$APP_NAME starting ..."
            nohup java ${JAVA_OPTS} -cp ${CLASSPATH_OPT} ${MAIN_CLASS} -f ${APP_HOME}/config/config.scala /tmp 2>> /dev/null >> /dev/null &
                        echo $! > $PID_PATH_NAME
            echo "$APP_NAME started ..."
        else
            echo "$APP_NAME is not running ..."
        fi
    ;;
esac 